# salt/root/files/concourse.service
[Unit]
Description=Concourse CI Server

[Service]
Type=simple
ExecStart=/usr/bin/concourse_start.sh

[Install]
WantedBy=multi-user.target
